# files/system/usr/etc/profile.d/ublue-os-just.sh
#I don't need a script to add my justfile to my home, and it's interfering with my current justfile, so this is just an empty file to prevent it from running.
