# files/auto-update-reboot/usr/lib/systemd/system/auto-update-reboot.timer
[Unit]
Description=Scheduled reboot with pending updates service

[Timer]
OnCalendar=*-*-* 5:30:00

[Install]
WantedBy=timers.target