# files/auto-update-reboot/usr/lib/systemd/system/auto-update-reboot.service
[Unit]
Description=Scheduled reboot with pending updates service

[Service]
Type=oneshot
ExecStart=/usr/bin/auto-update-reboot