# files/topgrade-updates/usr/lib/systemd/user/topgrade-update.service
[Unit]
Description=Auto Update With Topgrade Oneshot Service

[Service]
Type=oneshot
ExecStart=/usr/bin/topgrade-update