# config/files/usr/etc/profile.d/auraos-envs.sh
export XCURSOR_SIZE=36
export XCURSOR_THEME=Posys_cursor_Linux_black
export VISUAL=micro
export EDITOR=micro
export PAGER=moar
export SSH_AUTH_SOCK=$XDG_RUNTIME_DIR/ssh-agent.socket
export GTK_USE_PORTAL=1
